# packaging/systemd/conntrack.service
[Unit]
Description=Network Connection Tracker (eBPF)
After=network.target
Documentation=https://github.com/vponomarev/network-monitor

[Service]
Type=simple
ExecStart=/usr/local/bin/conntrack --config /etc/conntrack/config.yaml
Restart=on-failure
RestartSec=5

# Capabilities required for eBPF
AmbientCapabilities=CAP_BPF CAP_PERFMON CAP_NET_RAW CAP_SYS_ADMIN
CapabilityBoundingSet=CAP_BPF CAP_PERFMON CAP_NET_RAW CAP_SYS_ADMIN

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
